# gems/cronedit-0.3.0/test/examples/example1.cron
##__echo1__
* * * * * echo 1 >> /tmp/echo.log
##__echo2__
* * * * * echo 2 >> /tmp/echo.log
